# docker/supervisord.conf
# GameBox — Supervisor config
# Gestiona los procesos dentro del contenedor

[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:dbus]
command=dbus-daemon --system --fork
autostart=true
autorestart=false
startsecs=0

[program:pipewire]
command=pipewire
user=steam
environment=XDG_RUNTIME_DIR="/run/user/1000"
autostart=true
autorestart=true
startsecs=2
stopsignal=TERM

[program:pipewire-pulse]
command=pipewire-pulse
user=steam
environment=XDG_RUNTIME_DIR="/run/user/1000"
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM

[program:wireplumber]
command=wireplumber
user=steam
environment=XDG_RUNTIME_DIR="/run/user/1000"
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM

[program:sunshine]
command=sunshine --config /home/steam/.config/sunshine/sunshine.conf
user=steam
environment=XDG_RUNTIME_DIR="/run/user/1000"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/tmp/sunshine.log
stderr_logfile=/tmp/sunshine.err.log

[program:gamescope]
command=gamescope --headless --steam -- steam -steamos -pipewire -fulldesktopres -gamepadui
user=steam
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
